# fastdreg/files/etc/hotplug.d/iface/50-fastdstart
#!/bin/sh 
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
	sleep 3
	/lib/ffol/fastdreg/fastdreg.sh
}
